# drivers/net/wireless/ec19d/nanoradio-coredump.rules
ENV{NRX_EVENT}=="core",RUN+="/bin/cp -p $sys/$devpath/$env{NRX_CORENAME} /tmp"
